# birdwatch-stop.service
[Unit]
Description=Birdwatcher stopper

[Service]
Type=oneshot
ExecStart=/bin/true
