# systemd/rpi-reboot.service
[Unit]
Description=RaspberryPi fixer

[Service]
Type=oneshot
ExecStart=/usr/bin/reboot
